# config/files/usr/lib/systemd/system/coldturkey-enforce.service
[Unit]
Description=coldturkey - enforce NSFW blocking layers

[Service]
Type=oneshot
ExecStart=/usr/bin/env coldturkey enforce --quiet
