# config/files/usr/lib/systemd/system/coldturkey-enforce.timer
[Unit]
Description=coldturkey - periodic enforcement

[Timer]
OnBootSec=90
OnUnitActiveSec=5min
AccuracySec=30s

[Install]
WantedBy=timers.target
